# modules/podman/podman.socket
[Socket]
ListenStream=%t/podman/podman.sock
SocketMode=0666
SocketUser=root
SocketGroup=adm

[Install]
WantedBy=sockets.target
